# files/rstudio_server.in
#!/bin/sh

. /etc/rc.subr

name="rstudio_server"
rcvar=rstudio_server_enable

command="%%PREFIX%%/bin/rserver"

config_file="%%PREFIX%%/etc/rstudio/rserver.conf"
rsession_config_file="%%PREFIX%%/etc/rstudio/rsession.conf"
command_args="--config-file ${config_file} --rsession-config-file ${rsession_config_file}"

PATH=/sbin:/bin:/usr/sbin:/usr/bin:%%PREFIX%%/sbin:%%PREFIX%%/bin
export PATH

load_rc_config "$name"
run_rc_command  "$1"
